# packages/systemd/prady-swarm-coordinator.service
[Unit]
Description=Prady Swarm Coordinator Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/prady/platform/swarm-coordinator
ExecStart=/usr/local/bin/uvicorn swarm_coordinator:app --host 0.0.0.0 --port 8104
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
