# packages/systemd/prady-computer-use.service
[Unit]
Description=Prady Computer Use Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/prady/platform/computer-use
ExecStart=/usr/local/bin/uvicorn computer_use_service:app --host 0.0.0.0 --port 8106
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
